# etc/init.d/plymouth-openrc
#!/sbin/openrc-run
# Copyright (c) 2007-2021 The Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

description="Manages the Plymouth boot splash screen"
command="/usr/sbin/plymouthd"
command_background="yes"
pidfile="/run/plymouth/pid"

depend() {
	# Plymouth needs /dev and udev to be ready to interact with the framebuffer.
	need dev
	need udev
	# '-shutdown' ensures this service is not stopped during shutdown,
	# allowing its custom stop() logic to control the screen instead.
	keyword -shutdown
}

start_pre() {
	# Check kernel command line to see if splash is disabled.
	if grep -q -w 'nosplash\|plymouth.enable=0' /proc/cmdline; then
		einfo "Plymouth splash disabled via kernel command line."
		# Return a non-zero status to prevent the start() action.
		return 1
	fi
}

stop() {
	# This is not a standard stop action; it tells Plymouth to quit.
	# The 'plymouth quit' command handles a graceful exit.
	ebegin "Quitting Plymouth"
	plymouth quit
	eend $?
}
