# etc/init.d/plymouth-openrc
#!/sbin/openrc-run

description="Manages the Plymouth boot splash screen"
command="/usr/sbin/plymouthd"
command_args="--no-daemon"
supervisor=supervise-daemon

depend() {
	need dev udev
	keyword -shutdown
}

start_pre() {
	if grep -qwE 'nosplash|plymouth.enable=0' /proc/cmdline; then
		einfo "Plymouth splash disabled via kernel command line."
		return 1
	fi
}

stop() {
	ebegin "Quitting Plymouth"
	plymouth quit
	eend $?
}
